# systemd/user/random-wallpaper.service
[Unit]
Description=Set a random background image

[Service]
Type=simple
ExecStart=%h/.cargo/bin/random-wallpaper

[Install]
WantedBy=default.target
